# example/sources/default/templates/php-apache-web-rhel8/supervisord/php-fpm.conf
[program:php-fpm]
    command=/usr/sbin/php-fpm -FO
    stdout_logfile_maxbytes=0
    stdout_logfile=/dev/stdout
    stderr_logfile=/dev/stderr
    stderr_logfile_maxbytes=0
